# deploy/wedding-rsvp.service
[Unit]
Description=Wedding RSVP API (Maria Cristina & Diego)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/wedding-rsvp
ExecStart=/var/www/wedding-rsvp/wedding-rsvp
Restart=always
RestartSec=3

EnvironmentFile=/etc/wedding-rsvp.env

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/wedding-rsvp/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
